# packaging/DEBIAN/postinst
#!/bin/bash
set -e

echo "Configuring Punto Switcher..."

# 1. Setup Udev rules for uinput
echo 'KERNEL=="uinput", SUBSYSTEM=="misc", OPTIONS+="static_node=uinput", TAG+="uaccess", GROUP="input", MODE="0660"' > /etc/udev/rules.d/99-punto-uinput.rules
udevadm control --reload-rules
udevadm trigger || true

# 2. Add human users to input group
# We iterate users with UID >= 1000
for user in $(awk -F: '$3 >= 1000 && $1 != "nobody" {print $1}' /etc/passwd); do
    usermod -aG input "$user"
    echo "Added user $user to input group"
done

# 3. Setup Python Environment
# We assume code is in /opt/punto/app
APP_DIR="/opt/punto"

if [ ! -d "$APP_DIR/venv" ]; then
    echo "Creating virtual environment in $APP_DIR/venv..."
    python3 -m venv --system-site-packages "$APP_DIR/venv"
fi

# Install pure-python dependencies that might not be in apt
echo "Installing Python dependencies..."
"$APP_DIR/venv/bin/pip" install --upgrade pip
"$APP_DIR/venv/bin/pip" install num2words evdev

# 4. Fix permissions
chown -R root:root "$APP_DIR"
chmod -R 755 "$APP_DIR"

# 5. Enable user service globally
# This ensures it starts for all users on login
systemctl --global enable punto.service || true

echo "Punto Switcher installed successfully."
echo "Please LOGOUT and LOGIN again to apply group changes!"
