# DEBIAN/postinst
#!/bin/bash
set -e

# Создание директорий
mkdir -p /etc/punto
mkdir -p /etc/interception

# Определяем версию пакета
PKG_VERSION="2.7.3"
CONFIG_FILE="/etc/punto/config.yaml"
NEW_CONFIG="/usr/share/punto-switcher/config.yaml"

# Функция для получения версии из конфига
get_config_version() {
    if [ -f "$1" ]; then
        # Ищем строку вида "# Config version: X.X.X" в первых 10 строках
        grep -m1 "# Config version:" "$1" 2>/dev/null | sed 's/.*: //' || echo "0.0.0"
    else
        echo "0.0.0"
    fi
}

# Проверяем нужно ли обновить конфиг
if [ -f "$CONFIG_FILE" ]; then
    CURRENT_VERSION=$(get_config_version "$CONFIG_FILE")
    NEW_VERSION=$(get_config_version "$NEW_CONFIG")
    
    echo "Текущая версия конфига: $CURRENT_VERSION"
    echo "Новая версия конфига: $NEW_VERSION"
    
    # Если конфиг старый или версии нет - обновляем
    if [ "$CURRENT_VERSION" != "$NEW_VERSION" ] || [ "$CURRENT_VERSION" = "0.0.0" ]; then
        BACKUP_FILE="/etc/punto/config.yaml.backup-$(date +%Y%m%d-%H%M%S)"
        echo "Обновление конфигурации..."
        echo "  Backup: $BACKUP_FILE"
        cp "$CONFIG_FILE" "$BACKUP_FILE"
        cp "$NEW_CONFIG" "$CONFIG_FILE"
        echo "  ✓ Конфиг обновлён (старый сохранён в backup)"
    else
        echo "  Конфиг актуальный, обновление не требуется"
    fi
else
    # Первая установка
    echo "Установка начального конфига..."
    cp "$NEW_CONFIG" "$CONFIG_FILE"
fi

# Установка конфигурации udevmon
# Всегда обновляем, так как список клавиш может меняться
cp /usr/share/punto-switcher/udevmon.yaml /etc/interception/udevmon.yaml

# Включение и запуск сервиса
systemctl daemon-reload
systemctl enable udevmon || true
systemctl restart udevmon || true

echo "Punto Switcher for Ubuntu установлен успешно!"
echo "Author Anton Shalin <anton.shalin@gmail.com>"
echo "Используйте Pause для инверсии слова, Shift+Pause для инверсии выделенного текста."
