# packaging/punto.service
[Unit]
Description=Punto Switcher Daemon
After=graphical-session.target

[Service]
ExecStart=/opt/punto/venv/bin/python -m punto.daemon.main
Restart=always
RestartSec=3
Environment=PYTHONPATH=/opt/punto/app

[Install]
WantedBy=default.target
